# com.user.apple_music_lastfm_scrobbler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.apple_music_lastfm_scrobbler</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/apple_music_lastfm_scrobbler.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/apple_music_lastfm_scrobbler.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/apple_music_lastfm_scrobbler.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
<!--
INSTRUCTIONS:
- Replace "USERNAME" with your macOS username in the paths above
- Replace "/Users/USERNAME/apple_music_lastfm_scrobbler.sh" with the absolute path to your script if you placed it elsewhere
- This saves logs to /Users/USERNAME/Library/Logs/apple_music_lastfm_scrobbler.log
- Place this file at ~/Library/LaunchAgents/com.user.apple_music_lastfm_scrobbler.plist
- Load the service with: launchctl load ~/Library/LaunchAgents/com.user.apple_music_lastfm_scrobbler.plist
- Unload the service with: launchctl unload ~/Library/LaunchAgents/com.user.apple_music_lastfm_scrobbler.plist
-->

