# com.user.apple_music_lastfm_scrobbler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.apple_music_lastfm_scrobbler</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/PATH/TO/apple_music_lastfm_scrobbler.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/PATH/TO/REPO</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/Applications/Music.app</key>
            <true/>
        </dict>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/apple_music_lastfm_scrobbler.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/apple_music_lastfm_scrobbler.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
<!--
INSTRUCTIONS:
1. Copy this file to ~/Library/LaunchAgents/com.user.apple_music_lastfm_scrobbler.plist
2. Edit the file to replace:
   - USERNAME with your macOS username
   - PATH/TO/apple_music_lastfm_scrobbler.sh with the path to the script
   - PATH/TO/REPO with the path to the repository directory
3. Load the service: launchctl load ~/Library/LaunchAgents/com.user.apple_music_lastfm_scrobbler.plist
4. Unload the service: launchctl unload ~/Library/LaunchAgents/com.user.apple_music_lastfm_scrobbler.plist
5. Check logs at: ~/Library/Logs/apple_music_lastfm_scrobbler.log
-->

