# com.user.apple_music_lastfm_scrobbler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.apple_music_lastfm_scrobbler</string>
    <key>ProgramArguments</key>
    <array>
        <string>bash</string>
        <string>-c</string>
        <!-- This script will try multiple possible locations for the scrobbler script -->
        <string>
        # Try the default location in HOME directory first
        if [ -f "$HOME/apple_music_lastfm_scrobbler.sh" ]; then
            "$HOME/apple_music_lastfm_scrobbler.sh"
        # Then try the repository location if user kept it there
        elif [ -f "$HOME/Documents/xcode/applemusic_macos_scrobbler/apple_music_lastfm_scrobbler.sh" ]; then
            "$HOME/Documents/xcode/applemusic_macos_scrobbler/apple_music_lastfm_scrobbler.sh"
        # Allow for custom locations based on where the user cloned the repo
        elif [ -d "$HOME/github/applemusic_macos_scrobbler" ]; then
            "$HOME/github/applemusic_macos_scrobbler/apple_music_lastfm_scrobbler.sh"
        else
            echo "Error: Could not find the scrobbler script. Please update the path in $HOME/Library/LaunchAgents/com.user.apple_music_lastfm_scrobbler.plist"
            exit 1
        fi
        </string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>~/Library/Logs/apple_music_lastfm_scrobbler.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/apple_music_lastfm_scrobbler.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>

